# images/mysql/config/etc/cron.d/incremental-backup
*/5 * * * * root /usr/local/bin/incremental-backup.sh
